# trello-helper.timer
[Unit]
Description=Run foo weekly

[Timer]
OnCalendar=MON *-*-* 00:00:00
Persistent=true     
 
[Install]
WantedBy=timers.target
